# configs/camp/2018/gateway/etc/cron.d/dhcp_lease_influx
#

* * * * *  root  /opt/dhcp_influx/run.sh &> /dev/null
