# configs/camp/2018/monitoring_vm_ogling/etc/cron.d/aruba
*/1 * * * * root /root/aruba/aruba_cli_stats.sh > /dev/null
*/1 * * * * root /root/aruba/aruba_snmp_stats.sh > /dev/null
