# configs/camp/2018/monitoring_vm_ogling/etc/cron.d/update_gallery


*/10 * * * *  root  /opt/gallery/update.sh &> /dev/null
